# notification/crontab
*/1 * * * * /var/www/html/lodge-visa/notification/run.sh >> /var/www/html/lodge-visa/notification/log.txt 2>&1